# freeswitch/supervisord/supervisord.conf
[supervisord]
nodaemon=true

[program:freeswitch]
command=bash /startup.sh

[program:consul_agent]
command=/start_consul_agent.sh
stopsignal=INT
stopwaitsecs=5
stopasgroup=true

