# freeswitch/supervisord/supervisord.conf
[supervisord]
nodaemon=true

[program:freeswitch]
#command=/usr/bin/freeswitch -u freeswitch -g freeswitch -ncwait
command=bash /startup.sh
